# deploy/configs/bagisto-workers.conf
; Supervisor config for Bagisto workers
; Managed by Bagisto deploy scripts — do not edit manually

; ── Queue workers ──────────────────────────────────────────────────────────────
[program:bagisto-worker]
command=/usr/bin/php8.3 __APP_DIR__/artisan queue:work redis
    --queue=default,broadcastable
    --sleep=3
    --tries=3
    --max-time=3600
    --timeout=300
    --memory=256

process_name=%(program_name)s_%(process_num)02d
numprocs=__QUEUE_PROCESSES__
directory=__APP_DIR__

autostart=true
autorestart=true
startretries=10
startsecs=1

user=www-data
redirect_stderr=true
stdout_logfile=/var/log/bagisto/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=300

environment=APP_ENV="production",APP_DEBUG="false"

; ── Laravel Reverb (WebSockets) ────────────────────────────────────────────────
[program:bagisto-reverb]
command=/usr/bin/php8.3 __APP_DIR__/artisan reverb:start
    --host=0.0.0.0
    --port=__REVERB_PORT__
    --no-interaction

process_name=%(program_name)s
directory=__APP_DIR__

autostart=true
autorestart=true
startretries=10
startsecs=2

user=www-data
redirect_stderr=true
stdout_logfile=/var/log/bagisto/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=30

environment=APP_ENV="production",APP_DEBUG="false"

; ── Scheduler (Laravel artisan schedule:work) ──────────────────────────────────
[program:bagisto-scheduler]
command=/usr/bin/php8.3 __APP_DIR__/artisan schedule:work

process_name=%(program_name)s
directory=__APP_DIR__

autostart=true
autorestart=true
startretries=5
startsecs=1

user=www-data
redirect_stderr=true
stdout_logfile=/var/log/bagisto/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2

environment=APP_ENV="production",APP_DEBUG="false"
